# custom-modules/direct_osm2pgsql_update/files/init.conf
start on started rc
stop on stopping rc

setuid paladin
setgid nogroup


exec /opt/osm2pg_direct/osm2pg_direct.sh