# custom-modules/tiles/files/renderd.init.conf
description "Mapnik rendering daemon"
start on started rc
stop on stopping rc

setuid tiles
setgid nogroup

respawn
respawn limit 10 5

exec /usr/bin/renderd -f